# linux/airpods.service.in
[Unit]
Description=AirPods Daemon
After=bluetooth.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/airpods daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
